# deploy/debian/systemd/dream-crates-api.service
[Unit]
Description=Dream Crates API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dream-crates
Group=dream-crates
WorkingDirectory=/opt/dream-crates/app/backend
EnvironmentFile=/etc/dream-crates.env
ExecStart=/opt/dream-crates/app/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
